# bot.conf
[program:lectionbot]
directory=/root/LectionBot
command=python -m bot

autostart=true
autorestart=true

stderr_logfile=/root/logfile_err.log
stdout_logfile=/root/logfile.log